# deploy/systemd/memento-backup.service
[Unit]
Description=Memento local backup

[Service]
Type=oneshot
User=memento
Group=memento
ExecStart=/opt/memento/.venv/bin/memento --config /etc/memento/config.json backup --output /var/lib/memento/backups/latest
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/memento
