# deploy/systemd/memento-audit.service
[Unit]
Description=Memento repository audit

[Service]
Type=oneshot
User=memento
Group=memento
ExecStart=/opt/memento/.venv/bin/memento --config /etc/memento/config.json audit
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/memento
