# deploy/systemd/memento-backup.timer
[Unit]
Description=Run Memento backup daily (example disabled by default; maintenance commands require exclusive lease)

[Timer]
OnCalendar=daily
Persistent=true
Unit=memento-backup.service

# Intentionally no [Install] section: enable only after adding safe stop/conflict orchestration for memento.service.
